# valheim-backup.timer
[Unit]
Description=Valheim server backup timer

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
